# com.netmanager.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.netmanager.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <!-- IMPORTANT: EDIT THIS LINE to the full path of your start_net_manager.sh script -->
        <string>/path/to/your/network-manager-web/start_net_manager.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true />
</dict>
</plist>
